# hotline.service
[Unit]
Description=Hotline Client

[Service]
ExecStart=/root/hotline-client/startup.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
